# nfs/service/growfs@.timer
[Unit]
Description=Periodically run Grow File System on %f

[Timer]
OnActiveSec=0
OnCalendar=minutely
RandomizedDelaySec=30
