# nfs/service/growfs@.service
[Unit]
Description=Grow File System on %f

[Service]
Type=oneshot
ExecStart=/usr/bin/nfs-helper growfs %f
